# launchagents/com.example.codex-obsidian-save.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.example.codex-obsidian-save</string>

  <key>ProgramArguments</key>
  <array>
    <string>/path/to/node</string>
    <string>/path/to/codex-obsidian-save.js</string>
    <string>--background-save</string>
  </array>

  <key>WatchPaths</key>
  <array>
    <string>/path/to/home/.codex/session_index.jsonl</string>
  </array>

  <key>RunAtLoad</key>
  <false/>

  <key>ThrottleInterval</key>
  <integer>10</integer>

  <key>StandardOutPath</key>
  <string>/path/to/tmp/codex-obsidian-launchd.out</string>

  <key>StandardErrorPath</key>
  <string>/path/to/tmp/codex-obsidian-launchd.err</string>
</dict>
</plist>
